# scripts/com.teacher.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.teacher.backup</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/fredriksafsten/Projects/teacher/scripts/backup-prod.sh</string>
    </array>

    <!-- Run on load (login/wake) -->
    <key>RunAtLoad</key>
    <true/>

    <!-- Also run every 4 hours while awake as fallback -->
    <key>StartInterval</key>
    <integer>14400</integer>

    <!-- Environment variables -->
    <key>EnvironmentVariables</key>
    <dict>
        <key>PATH</key>
        <string>/opt/homebrew/bin:/usr/local/bin:/usr/bin:/bin</string>
    </dict>

    <!-- Logging -->
    <key>StandardOutPath</key>
    <string>/Users/fredriksafsten/Projects/teacher/backups/launchd.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/fredriksafsten/Projects/teacher/backups/launchd-error.log</string>
</dict>
</plist>
